# agenticops-frontend.service
[Unit]
Description=AgenticOps Frontend Service (Vue)
After=network.target agenticops-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/agenticops/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/env npm run preview -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=10
StandardOutput=append:/path/to/agenticops/logs/frontend.log
StandardError=append:/path/to/agenticops/logs/frontend.log

[Install]
WantedBy=multi-user.target
