# netops-backend.service
[Unit]
Description=NetOps AI Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netops/backend
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /opt/netops/backend/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target