# netops-backend.service
[Unit]
Description=AgenticOps Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/netops_bs/backend
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/netops_bs/backend/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
